# rootfs/etc/s6-overlay/s6-rc.d/init-adduser/run
#! /usr/bin/with-contenv oi
# shellcheck shell=bash
declare ARG_UID="${_UID:-568}";
declare ARG_GID="${_GID:-568}";
declare ARG_USER="${_USER:-admin}";

function init::chezmoi {
  # Initialize chezmoi or check of updates if already initialized.
  #
  local repo="${1:-}";  # repo that chezmoi will checkout and apply
  local user_name="${2:-admin}"  # name of non-root user

  if oi::var.is_empty "${repo}"; then
    oi::log.info "git repo for chezmoi managed dotfiles not provided";
    return
  fi

  if oi::fs.directory_exists "/root/.local/share/chezmoi"; then
    chezmoi update --apply --color true --force \
      || oi::log.warning "Failed to update & apply chezmoi dotfiles";
  else
    chezmoi init "${repo}" --apply --color true \
      || oi::exit.error "Failed to initialize & apply chezmoi dotfiles";
  fi
}

function init::user {
  # Initialize user.
  #
  local user_name="${1:-admin}";
  local arg_gid="${2:-568}";
  local arg_uid="${3:-568}";

  groupmod --non-unique --gid "$arg_gid" admin \
    || oi::exit.error "Failed to update admin group GID";

  usermod --non-unique --uid "$arg_uid" admin \
    || oi::exit.error "Failed to update admin UID";

  if [[ "${ARG_USER}" != "admin" ]]; then
    usermod --home "/home/${ARG_USER}" --login "${ARG_USER}" --move-home --non-unique admin \
      || oi::exit.error "Failed to update admin user to ${ARG_USER}";
  fi

  for sub_dir in 'bin' 'lib' 'share'; do
    if ! oi::fs.directory_exists "/home/${user_name}/.local/${sub_dir}"; then
      mkdir -p "/home/${user_name}/.local/${sub_dir}" \
        || oi::exit.error "Failed to create directory '/home/${user_name}/.local/${sub_dir}'";
      lsiown "${user_name}:${user_name}" "/home/${user_name}/.local" \
        || oi::exit.error "Failed to change ownership of '/home/${user_name}/.local/${sub_dir}'";
    fi
  done

  lsiown "${user_name}:${user_name}" /app;
  lsiown "${user_name}:${user_name}" /config;
  lsiown "${user_name}:${user_name}" /data;
  lsiown "${user_name}:${user_name}" /defaults;
}

init::user "${ARG_USER}" "${ARG_GID}" "${ARG_UID}";
init::chezmoi "${CHEZMOI_REPO:-}";
