# rootfs/etc/s6-overlay/s6-rc.d/init-adduser/run
#! /usr/bin/with-contenv oi
# shellcheck shell=bash

ARG_UID="${_UID:-568}";
ARG_GID="${_GID:-568}";
ARG_USER="${_USER:-admin}";

groupmod --non-unique --gid "$ARG_GID" admin \
  || oi::exit.error "Failed to update admin group GID";

usermod --non-unique --uid "$ARG_UID" admin \
    || oi::exit.error "Failed to update admin UID";

if [[ "${ARG_USER}" != "admin" ]]; then
  usermod --home "/home/${ARG_USER}" --login "${ARG_USER}" --move-home --non-unique admin \
    || oi::exit.error "Failed to update admin user to ${ARG_USER}";
fi

lsiown "${ARG_USER}:${ARG_USER}" /app
lsiown "${ARG_USER}:${ARG_USER}" /config
lsiown "${ARG_USER}:${ARG_USER}" /data
lsiown "${ARG_USER}:${ARG_USER}" /defaults
